# rootfs/etc/services.d/nginx/run
#!/usr/bin/env bashio
# ==============================================================================
# Start NGINX service
# ==============================================================================


bashio::log.info "Starting Nginx..."
exec nginx